# conf/supervisord.ini
[program:flask_rest_display]
user=nginx
command=/var/www/flask_rest_display/env/bin/uwsgi --ini conf/uwsgi.ini
directory=/var/www/flask_rest_display
stdout_logfile=/var/log/supervisor/flask_rest_display.log
stderr_logfile=/var/log/supervisor/flask_rest_display.log

[program:flask_rest_display_celery]
user=nginx
command=/var/www/flask_rest_display/env/bin/celery -A celery_app.celery worker -B -Q input_q,retry_input_q,image_monitor_q -l info -c 5 -s /tmp/flask_rest_display_celery
directory=/var/www/consumer_financial_rc_apply_yidai
stdout_logfile=/var/log/supervisor/flask_rest_display_celery.log
stderr_logfile=/var/log/supervisor/flask_rest_display_celery.log
